# systemd/newsmon.service
[Unit]
Description=NewsMon FastAPI Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=maan
Group=maan
WorkingDirectory=/home/maan/newsmon/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/newsmon.env
ExecStart=/home/maan/newsmon/backend/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
